# support/supervisord/cert_manage.ini
[program:cms]
directory=/opt/cms
environment=PATH="/opt/py3/bin:%(ENV_PATH)s"
command=/opt/py3/bin/python run_server.py all
pidfile=/opt/cms/cms.pid
autostart=true
autorestart=true
stdout_logfile=/opt/cms/logs/cms.log
stderr_logfile=/opt/cms/logs/cms_error.log